# test/supervisord.conf
# Runs the ARK-OS userspace stack in one container, all on shared localhost — mirroring
# the real device, where every service binds 127.0.0.1 and nginx is the only public face.
# (system/service/connection-manager hardcode host 127.0.0.1, so they must co-reside with
# the gateway that proxies to them.) Each program's stdout is forwarded to the container
# stdout so `docker compose logs -f arkos` shows everything.
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:ark-ui-backend]
command=/usr/local/bin/node /usr/lib/ark-os/ark-ui-backend/index.js
environment=PORT="3000",NODE_ENV="production"
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:system-manager]
command=python3 system_manager.py
directory=/usr/lib/ark-os/python/system-manager
environment=PORT="3004"
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:service-manager]
command=python3 service_manager.py
directory=/usr/lib/ark-os/python/service-manager
environment=PORT="3002"
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:connection-manager]
command=python3 connection_manager.py
directory=/usr/lib/ark-os/python/connection-manager
environment=PORT="3001"
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:autopilot-manager]
command=python3 autopilot_manager.py --host 127.0.0.1 --port 3003 --connection-string udpin:127.0.0.1:14571
directory=/usr/lib/ark-os/python/autopilot-manager
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:mavlink-stub]
command=python3 /usr/lib/ark-os/test/mavlink_stub.py
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
