# packaging/service-files/jetson/autopilot-manager.service
[Unit]
Description=Microservice backend for autopilot mavlink interactions
Wants=network.target network-online.target
After=network-online.target

[Service]
Type=simple
User=jetson
Group=jetson
ExecStart=/usr/lib/ark-os/venv/bin/python3 /usr/lib/ark-os/python/autopilot_manager.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=PORT=3003
# clock_settime: adopt the FC's GNSS time when the local clock has no NTP sync.
AmbientCapabilities=CAP_SYS_TIME

[Install]
WantedBy=multi-user.target ark-os.target
